# 2025-11-02_r32350-509c254c48_next-r4.6.14.rss.mtk/overlay-files/etc/uci-defaults/99-remove-apk-message
#!/bin/sh

# Remove the APK migration message from /etc/profile

PROFILE="/etc/profile"

if [ -f "$PROFILE" ]; then
    # Remove the entire APK message block
    sed -i '/^if \[ -x \/usr\/bin\/apk \]; then$/,/^fi$/d' "$PROFILE"
fi

exit 0