# files/etc/profile.d/zsh.sh
# by ogglord
if [[ $- == *i* ]]; then
    if [ -x /usr/bin/zsh ]; then
    	exec /usr/bin/zsh
    fi
fi
