# 2025-11-05_r32350-509c254c48_next-r4.6.14.rss.mtk/overlay-files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
#
# Note: Network tweaks (tc qdisc, ethtool) are now handled by
# /etc/init.d/network-tweaks service which runs on each boot
# with proper error handling.

exit 0